# toolchain/eldk-5.2.1/armv5te/sysroots/armv5te-linux-gnueabi/var/lib/opkg/info/consolekit.list
/usr/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Seat.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Manager.xml
/usr/share/dbus-1/interfaces/
/lib/systemd/
/var/log/
/lib/systemd/system/
/etc/ConsoleKit/
/usr/lib/ConsoleKit/run-session.d/
/usr/lib/ConsoleKit/
/lib/systemd/system/console-kit-daemon.service
/usr/lib/libck-connector.so.0
/usr/share/polkit-1/actions/org.freedesktop.consolekit.policy
/var/volatile/log/ConsoleKit/
/usr/lib/ConsoleKit/scripts/
/usr/libexec/ck-collect-session-info
/lib/systemd/system/kexec.target.wants/
/usr/libexec/ck-get-x11-server-pid
/usr/share/dbus-1/system-services/org.freedesktop.ConsoleKit.service
/lib/systemd/system/basic.target.wants/console-kit-log-system-start.service
/etc/ConsoleKit/seats.d/00-primary.seat
/etc/ConsoleKit/run-session.d/
/usr/lib/ConsoleKit/scripts/ck-system-stop
/usr/sbin/ck-log-system-restart
/usr/libexec/ck-get-x11-display-device
/lib/systemd/system/reboot.target.wants/console-kit-log-system-restart.service
/etc/ConsoleKit/run-seat.d/
/lib/systemd/system/basic.target.wants/
/usr/share/polkit-1/
/var/log/ConsoleKit/
/lib/systemd/system/poweroff.target.wants/
/lib/systemd/system/console-kit-log-system-stop.service
/lib/systemd/system/poweroff.target.wants/console-kit-log-system-stop.service
/lib/systemd/system/halt.target.wants/
/etc/ConsoleKit/seats.d/
/usr/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Session.xml
/usr/sbin/ck-log-system-stop
/var/volatile/log/
/lib/systemd/system/kexec.target.wants/console-kit-log-system-restart.service
/etc/dbus-1/system.d/ConsoleKit.conf
/lib/systemd/system/console-kit-log-system-restart.service
/usr/sbin/ck-log-system-start
/lib/systemd/system/halt.target.wants/console-kit-log-system-stop.service
/usr/lib/ConsoleKit/scripts/ck-system-restart
/usr/lib/ConsoleKit/run-seat.d/
/usr/bin/ck-list-sessions
/usr/sbin/console-kit-daemon
/usr/lib/libck-connector.so.0.0.0
/usr/bin/ck-history
/lib/systemd/system/console-kit-log-system-start.service
/usr/share/polkit-1/actions/
/var/volatile/run/ConsoleKit/
/usr/bin/ck-launch-session
/lib/systemd/system/reboot.target.wants/
