# rootfs/picosafe_rootfs/var/lib/dpkg/info/avahi-daemon.list
/etc
/etc/default
/etc/default/avahi-daemon
/etc/init.d
/etc/init.d/avahi-daemon
/etc/resolvconf
/etc/resolvconf/update-libc.d
/etc/resolvconf/update-libc.d/avahi-daemon
/etc/avahi
/etc/avahi/hosts
/etc/avahi/avahi-daemon.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/avahi-dbus.conf
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/avahi-daemon
