# rootfs/picosafe_rootfs/var/lib/dpkg/info/clamav-freshclam.list
/.
/etc
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/clamav-freshclam-ifupdown
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/clamav-freshclam-ifupdown
/etc/network
/etc/network/if-down.d
/etc/network/if-down.d/clamav-freshclam-ifupdown
/etc/network/if-up.d
/etc/network/if-up.d/clamav-freshclam-ifupdown
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/clamav-freshclam
/etc/init.d
/etc/init.d/clamav-freshclam
/etc/clamav
/etc/clamav/onerrorexecute.d
/etc/clamav/onupdateexecute.d
/var
/var/lib
/var/lib/clamav
/usr
/usr/bin
/usr/bin/freshclam
/usr/share
/usr/share/doc
/usr/share/doc/clamav-freshclam
/usr/share/doc/clamav-freshclam/examples
/usr/share/doc/clamav-freshclam/examples/freshclam.conf
/usr/share/doc/clamav-freshclam/changelog.gz
/usr/share/doc/clamav-freshclam/mirror-list.gz
/usr/share/doc/clamav-freshclam/changelog.Debian.gz
/usr/share/doc/clamav-freshclam/copyright
/usr/share/bug
/usr/share/bug/clamav-freshclam
/usr/share/bug/clamav-freshclam/script
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/freshclam.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/freshclam.1.gz
/usr/share/doc/clamav-freshclam/README.gz
/usr/share/doc/clamav-freshclam/NEWS.Debian.gz
/usr/share/doc/clamav-freshclam/AUTHORS
/usr/share/doc/clamav-freshclam/README.Debian.gz
/usr/share/doc/clamav-freshclam/BUGS
/usr/share/doc/clamav-freshclam/FAQ
