# rootfs/picosafe_rootfs/var/lib/dpkg/info/cryptsetup.list
/.
/sbin
/sbin/cryptdisks_stop
/sbin/cryptsetup
/sbin/cryptdisks_start
/lib
/lib/cryptsetup
/lib/cryptsetup/checks
/lib/cryptsetup/checks/ext2
/lib/cryptsetup/checks/blkid
/lib/cryptsetup/checks/un_blkid
/lib/cryptsetup/checks/xfs
/lib/cryptsetup/checks/swap
/lib/cryptsetup/cryptdisks.functions
/lib/cryptsetup/scripts
/lib/cryptsetup/scripts/decrypt_ssl
/lib/cryptsetup/scripts/decrypt_gnupg
/lib/cryptsetup/scripts/decrypt_opensc
/lib/cryptsetup/scripts/decrypt_openct
/lib/cryptsetup/scripts/decrypt_keyctl
/lib/cryptsetup/scripts/decrypt_derived
/lib/cryptsetup/scripts/passdev
/lib/cryptsetup/askpass
/usr
/usr/sbin
/usr/sbin/luksformat
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/crypttab.5.gz
/usr/share/man/man8
/usr/share/man/man8/cryptdisks_stop.8.gz
/usr/share/man/man8/cryptsetup.8.gz
/usr/share/man/man8/luksformat.8.gz
/usr/share/man/man8/cryptdisks_start.8.gz
/usr/share/locale
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/vi
/usr/share/locale/vi/LC_MESSAGES
/usr/share/locale/vi/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/cryptsetup.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/cryptsetup.mo
/usr/share/doc
/usr/share/doc/cryptsetup
/usr/share/doc/cryptsetup/AUTHORS
/usr/share/doc/cryptsetup/changelog.Debian.gz
/usr/share/doc/cryptsetup/examples
/usr/share/doc/cryptsetup/examples/gen-ssl-key
/usr/share/doc/cryptsetup/changelog.gz
/usr/share/doc/cryptsetup/NEWS.Debian.gz
/usr/share/doc/cryptsetup/README.openct.gz
/usr/share/doc/cryptsetup/README.initramfs.gz
/usr/share/doc/cryptsetup/README.Debian.gz
/usr/share/doc/cryptsetup/README.gnupg
/usr/share/doc/cryptsetup/README.remote.gz
/usr/share/doc/cryptsetup/README.opensc.gz
/usr/share/doc/cryptsetup/copyright
/usr/share/doc/cryptsetup/README.keyctl
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/cryptsetup
/usr/share/bug
/usr/share/bug/cryptsetup
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/local-bottom
/usr/share/initramfs-tools/scripts/local-bottom/cryptopensc
/usr/share/initramfs-tools/scripts/local-top
/usr/share/initramfs-tools/scripts/local-top/cryptopensc
/usr/share/initramfs-tools/scripts/local-top/cryptroot
/usr/share/initramfs-tools/conf-hooks.d
/usr/share/initramfs-tools/conf-hooks.d/cryptsetup
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/cryptopenct
/usr/share/initramfs-tools/hooks/cryptpassdev
/usr/share/initramfs-tools/hooks/cryptgnupg
/usr/share/initramfs-tools/hooks/cryptopensc
/usr/share/initramfs-tools/hooks/cryptroot
/usr/share/initramfs-tools/hooks/cryptkeyctl
/etc
/etc/default
/etc/default/cryptdisks
/etc/init.d
/etc/init.d/cryptdisks
/etc/init.d/cryptdisks-early
/usr/sbin/cryptdisks_stop
/usr/sbin/cryptdisks_start
