# rootfs/picosafe_rootfs/var/lib/dpkg/info/hdparm.conffiles
/etc/udev/hdparm.rules
/etc/init.d/hdparm
/etc/hdparm.conf
/etc/default/hdparm
/etc/apm/event.d/20hdparm
