# rootfs/picosafe_rootfs/var/lib/dpkg/info/nscd.list
/etc
/etc/init.d
/etc/init.d/nscd
/etc/nscd.conf
/var/cache/nscd
