# rootfs/picosafe_rootfs/var/lib/dpkg/info/nscd.postrm
#!/bin/sh
set -e

case "$1" in
    purge)
	rm -rf /var/cache/nscd
	;;
    *)
    	;;
esac

# Automatically added by dh_installinit
if [ "$1" = "purge" ] ; then
	update-rc.d nscd remove >/dev/null
fi
# End automatically added section


exit 0
