# rootfs/picosafe_rootfs/var/lib/dpkg/info/procps.conffiles
/etc/sysctl.conf
/etc/sysctl.d/README.sysctl
/etc/init.d/procps
