# rootfs/picosafe_rootfs/var/lib/dpkg/info/tmux.postrm
#!/bin/sh
set -e

if [ "$1" = purge ] || [ "$1" = remove ]; then
  rm -rf /var/run/tmux
fi

# Automatically added by dh_installinit
if [ "$1" = "purge" ] ; then
	update-rc.d tmux-cleanup remove >/dev/null
fi
# End automatically added section

